# scripts/autostart/virtualbox.desktop
[Desktop Entry]
Name=AutoStartVirtualBox
Exec=/home/user/autostart-virtualbox
Terminal=false
Type=Application
StartupNotify=true
